# image/mbr/simple_dual/linko/launcher/linko_launcher.service
[Unit]
Description=Kiosk Linko Launcher
After=network.target

[Service]
ExecStart=/usr/bin/java -jar /opt/kiosk/LinkoLauncher.jar
User=player
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target